# tinyWeather.timer
[Unit]
Description=Daily weather alerts

[Timer]
OnCalendar=daily
AccuracySec=23h
Persistent=true

[Install]
WantedBy=default.target
